# docker/build/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:mariadb]
command=/usr/sbin/mysqld --user=root
autostart=true
autorestart=true
stderr_logfile=/var/log/mariadb.err.log
stdout_logfile=/var/log/mariadb.out.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
redirect_stderr=true

[program:node]
command=node server
